# supervisor/comfyui.conf
[program:comfyui]
command=/bin/bash -lc 'cd /workspace/ComfyUI && /venv/main/bin/python main.py --listen 0.0.0.0 --port 18188'
autostart=true
autorestart=true
startsecs=10
stdout_logfile=/tmp/comfyui.log
stderr_logfile=/tmp/comfyui.err.log
stdout_logfile_maxbytes=50MB
stderr_logfile_maxbytes=50MB
